# launchd/com.argos.douban.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  阿尔戈斯计划 · 豆瓣每日抓取 launchd 任务
  每天 10:30 运行 scripts/run-daily.mjs,日志写到 logs/。

  ⚠️ 安装前必须替换下面两处占位:
    1. 三个 /替换为你的实际路径/  →  项目绝对路径,例如 /Users/你的用户名/douban-cloud
    2. /替换为node路径/node      →  运行 `which node` 得到,例如 /opt/homebrew/bin/node
       (Apple Silicon 一般是 /opt/homebrew/bin/node,Intel 一般是 /usr/local/bin/node)

  安装:  cp launchd/com.argos.douban.daily.plist ~/Library/LaunchAgents/
         launchctl load ~/Library/LaunchAgents/com.argos.douban.daily.plist
  卸载:  launchctl unload ~/Library/LaunchAgents/com.argos.douban.daily.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.argos.douban.daily</string>

  <!-- 运行命令:node scripts/run-daily.mjs -->
  <key>ProgramArguments</key>
  <array>
    <string>/替换为node路径/node</string>
    <string>/替换为你的实际路径/douban-cloud/scripts/run-daily.mjs</string>
  </array>

  <!-- 工作目录(让脚本里的相对路径正确) -->
  <key>WorkingDirectory</key>
  <string>/替换为你的实际路径/douban-cloud</string>

  <!-- 从 .env 加载环境变量较麻烦,这里推荐两种方式之一:
       A) 在下面 EnvironmentVariables 里直接写(简单,但明文在 plist)
       B) 让 run-daily 自己读 .env(本项目脚本通过 process.env;
          launchd 不自动读 .env,故仍需在此或在外部注入)。
       为简单起见,建议用 setup-macmini.sh 自动把 .env 内容注入到这里。 -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>SUPABASE_URL</key><string>__SUPABASE_URL__</string>
    <key>SUPABASE_SERVICE_ROLE_KEY</key><string>__SUPABASE_SERVICE_ROLE_KEY__</string>
    <key>SUPABASE_ANON_KEY</key><string>__SUPABASE_ANON_KEY__</string>
    <key>DOUBAN_COOKIE</key><string>__DOUBAN_COOKIE__</string>
    <key>DOUBAN_PROXY</key><string>__DOUBAN_PROXY__</string>
    <key>PATH</key><string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>

  <!-- 每天 10:30 触发 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>10</integer>
    <key>Minute</key><integer>30</integer>
  </dict>

  <!-- 日志 -->
  <key>StandardOutPath</key>
  <string>/替换为你的实际路径/douban-cloud/logs/daily.out.log</string>
  <key>StandardErrorPath</key>
  <string>/替换为你的实际路径/douban-cloud/logs/daily.err.log</string>

  <!-- 错过时间(如机器当时睡眠)上线后补跑一次 -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
